# packaging/lib/systemd/system/minifw-ai.service
[Unit]
Description=MiniFW-AI — Lightweight AI-Powered Firewall Companion
Documentation=file:///opt/minifw-ai/README.md
After=network.target ritapi-advanced.service
Wants=ritapi-advanced.service

[Service]
Type=simple
User=ritapi
Group=ritapi
WorkingDirectory=/opt/minifw-ai
EnvironmentFile=/etc/minifw-ai/env
ExecStart=/opt/minifw-ai/.venv/bin/python -m minifw_ai.main
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/log/minifw-ai /etc/minifw-ai

[Install]
WantedBy=multi-user.target
